# .config/systemd/user/autotiling.service
[Unit]
ConditionEnvironment=SWAYSOCK

[Service]
ExecStart=autotiling
Restart=on-failure

[Install]
